# config/supervisor/2-openbox.conf
[program:Openbox]
command=bash -c 'sleep 3 && openbox'
priority=2
autorestart=true
stderr_logfile=/var/log/openbox.stderr.log
stdout_logfile=/var/log/openbox.stdout.log
